# deploy/server.service
[Unit]
Description=Amazon Scraper v3 Server
After=network.target

[Service]
Type=simple
User=scraper
WorkingDirectory=/opt/amazon-scraper-v3
ExecStart=/opt/amazon-scraper-v3/.venv/bin/python run_server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
